# performance-rig.service
[Unit]
Description=Live Performance Rig Controller
After=network.target sound.target pipewire.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/rig
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
ExecStart=/home/pi/rig/venv/bin/python /home/pi/rig/controller.py
Restart=on-failure
RestartSec=5

# Audio settings
Environment="PIPEWIRE_LATENCY=512/48000"

[Install]
WantedBy=multi-user.target
